# guard.service
[Unit]
Description=Security guard

[Service]
Type=oneshot
ExecStart=/etc/systemd/system/guard.sh
