# hide_location.service
[Unit]
Description=Hide location DNS

[Service]
Type=oneshot
ExecStart=/etc/systemd/system/hide_location.sh
