# open_vpn.service
[Unit]
Description=Open VPN service

[Service]
Type=oneshot
ExecStart=/etc/systemd/system/open_vpn.sh
